# rootfs/etc/s6-overlay/s6-rc.d/postgres-post/run
#!/usr/bin/env bash
set -euo pipefail

source "${BASHIO_LIB:-/usr/lib/bashio/bashio.sh}"

readonly RUNTIME_DIR="/run/postgresql"
readonly READY_TIMEOUT_SECONDS=60

ready=0
for ((attempt=1; attempt<=READY_TIMEOUT_SECONDS; attempt++)); do
  if pg_isready -h "${RUNTIME_DIR}" -p 5432 -U postgres >/dev/null 2>&1; then
    ready=1
    break
  fi
  sleep 1
done

if [[ "${ready}" -ne 1 ]]; then
  bashio::log.fatal "PostgreSQL did not become ready within ${READY_TIMEOUT_SECONDS} seconds"
  exit 1
fi

/usr/local/bin/provisioning-entrypoint

bashio::log.info "Publishing PostgreSQL service information"
PAYLOAD=$(
  bashio::var.json \
    host "$(hostname)" \
    port "^5432"
)
if ! bashio::services.publish "postgresql" "${PAYLOAD}"; then
  bashio::log.warning "Failed to publish PostgreSQL service information"
fi

bashio::log.info "PostgreSQL 18 ready"
