# rootfs/etc/s6-overlay/s6-rc.d/postgres-pre/up
#!/usr/bin/env bash
set -euo pipefail

source /usr/lib/bashio/bashio

readonly DATA_DIR="/data/postgres"
readonly RUNTIME_DIR="/run/postgresql"
readonly LOG_DIR="/var/log/postgresql"

mkdir -p "${DATA_DIR}" "${RUNTIME_DIR}" "${LOG_DIR}"
chown -R postgres:postgres "${DATA_DIR}" "${RUNTIME_DIR}" "${LOG_DIR}"

if [[ ! -f "${DATA_DIR}/PG_VERSION" ]]; then
  bashio::log.info "Initializing PostgreSQL 18 cluster"
  runuser -u postgres -- /usr/lib/postgresql/18/bin/initdb \
    -D "${DATA_DIR}" \
    --username=postgres \
    --encoding=UTF8 \
    --locale=C.UTF-8 \
    --auth-local=trust \
    --auth-host=scram-sha-256
fi

/usr/bin/render-postgres-config
chown postgres:postgres "${DATA_DIR}/postgresql.conf" "${DATA_DIR}/pg_hba.conf" "${DATA_DIR}/pg_ident.conf"
